# utils/nightly/gmsh_linux_nightly.cron
#
# this will launch a new build every night at 1h01
#
1 1 * * * cd ${HOME}/src/gmsh/utils/nightly && /usr/local/bin/ctest -S gmsh_linux_nightly.ctest
