# utils/nightly/gmsh_mac_nightly.cron
#
# this will launch a new build every night at 0h01
#
1 0 * * * cd ${HOME}/src/gmsh/utils/nightly && /usr/local/bin/ctest -S gmsh_mac_nightly.ctest
